# supervisor/conf.d/gopro108.conf
[program:gopro108]
command=/projects/kiosk-app/backend/venv/bin/python /projects/kiosk-app/backend/GoProStream108.py
user=root
stdout_logfile=/projects/kiosk-app/logs/gopro108.out.log
stderr_logfile=/projects/kiosk-app/logs/gopro108.err.log
autostart=true
autorestart=true
